# dot_config/systemd/user/battery.service
[Unit]
Description=Checks battery level every 5 minutes

[Service]
Type=simple
ExecStart=%h/.local/bin/battery-warning.sh
